# home/dot_config/systemd/user/hyda-updater.service
[Unit]
Description=Runs Hyda Updater Every 6 Hours

[Service]
Type=oneshot
ExecStart=%h/.local/share/hyda/bin/hyda-updater.sh notify
